# aes-app.service
[Unit]
Description=AES Encryption Web Application
After=network.target

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/var/www/aes
Environment="PATH=/var/www/aes/venv/bin"
ExecStart=/var/www/aes/venv/bin/gunicorn -c gunicorn.conf.py wsgi:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=3

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/aes

[Install]
WantedBy=multi-user.target 